# wttr-locale.cron
SHELL=/bin/bash

#min  hour  day   month wday  command
#------------------------------------
*/15  *     *     *     *     /home/elektronik/1.sh
