# wttr-locale.cron
SHELL=/bin/bash

#min  hour  day   month wday  command
#------------------------------------
*/15  *     *     *     *     /<DIR>/wttr-locale.sh         # <------ enter location
